# sh/beacon.service
[Unit]
Description=
After=syslog.target
[Service]
User=beacon
ExecStart=/var/beacon/beacon.jar
[Install]
WantedBy=multi-user.target